# debian/postinst
#!/bin/bash
# postinst script for cryptnox-cli
set -e

case "$1" in
    configure)
        # Install Python dependencies from PyPI (not all available in Debian repos)
        echo "Installing Python dependencies..."

        # All dependencies from cryptnox-cli pyproject.toml
        DEPS="cryptnox-sdk-py>=1.0.2 lazy-import tabulate argparse appdirs base58 boto3 ecdsa colander pytz requests stdiomask web3"

        # Try with --break-system-packages first (Python 3.11+), fallback to without
        if pip3 install --break-system-packages --quiet $DEPS 2>/dev/null; then
            echo "Dependencies installed successfully"
        elif pip3 install --quiet $DEPS 2>/dev/null; then
            echo "Dependencies installed successfully"
        else
            echo "Warning: Some dependencies may not have installed correctly"
            echo "Try running: pip3 install cryptnox-cli"
        fi

        # Ensure pcscd is running
        if command -v systemctl &>/dev/null; then
            systemctl enable pcscd 2>/dev/null || true
            systemctl start pcscd 2>/dev/null || true
        fi
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
